# service/btrfs-balance.service
[Unit]
Description=btrfs balance service

[Service]
Type=simple
ExecStart=/usr/sbin/btrfs-balance

[Install]
WantedBy=btrfs-balancer.timer
